# configs/systemd/user/level-bridge@.service
[Unit]
Description=Level Bridge metering tap: %i
After=pipewire.service
Wants=pipewire.service

[Service]
Type=simple
EnvironmentFile=/etc/level-bridge/%i.env
ExecStart=/home/ela/bin/level-bridge \
    ${LINK_MODE} \
    --node-name ${NODE_NAME} \
    --mode ${MODE} \
    --target ${TARGET} \
    --levels-listen tcp:127.0.0.1:${PORT} \
    --channels ${CHANNELS} \
    --rate ${RATE}

Restart=on-failure
RestartSec=2

# RT scheduling: PipeWire's data loop thread runs at the server's
# priority (SCHED_FIFO/88). The binary itself runs at normal priority.
OOMScoreAdjust=-500
LimitMEMLOCK=infinity
LimitRTPRIO=88

# Security hardening — PipeWire client profile (SEC-PW-CLIENT)
# See F-035 for rationale.
NoNewPrivileges=true
PrivateTmp=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
CapabilityBoundingSet=
RestrictSUIDSGID=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged
ProtectSystem=full
ProtectHome=read-only

[Install]
WantedBy=default.target
