# configs/systemd/user/pcm-bridge@.service
[Unit]
Description=PCM Bridge audio stream: %i
After=pipewire.service
Wants=pipewire.service

[Service]
Type=simple
EnvironmentFile=/etc/pcm-bridge/%i.env
ExecStart=/home/ela/bin/pcm-bridge \
    --mode ${MODE} \
    --target ${TARGET} \
    --listen tcp:127.0.0.1:${PORT} \
    --channels ${CHANNELS} \
    --rate ${RATE} \
    --quantum ${QUANTUM} \
    --levels-listen ${LEVELS_LISTEN}

Restart=on-failure
RestartSec=2

# RT scheduling: inherited from PipeWire.
# pcm-bridge uses RT_PROCESS streams, so PipeWire's data thread invokes
# our process callbacks at its own RT priority (SCHED_FIFO/88). The binary
# itself runs at normal priority; only the PW callback thread is RT.
OOMScoreAdjust=-500
LimitMEMLOCK=infinity
LimitRTPRIO=88

# Security hardening — PipeWire client profile (SEC-PW-CLIENT)
# See F-035 for rationale.
NoNewPrivileges=true
PrivateTmp=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
CapabilityBoundingSet=
RestrictSUIDSGID=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged
ProtectSystem=full
ProtectHome=read-only

[Install]
WantedBy=default.target
