# configs/systemd/user/pi4-audio-webui.service
[Unit]
Description=Pi4 Audio Workstation monitoring web UI (D-020)
After=pipewire.service
Wants=pipewire.service

[Service]
Type=simple
WorkingDirectory=/home/ela/web-ui
ExecStart=/home/ela/audio-workstation-venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8080 --workers 1 --ssl-keyfile /home/ela/web-ui/key.pem --ssl-certfile /home/ela/web-ui/cert.pem
Environment=PI_AUDIO_MOCK=0
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=JACK_NO_START_SERVER=1
Environment=LD_LIBRARY_PATH=/usr/lib/aarch64-linux-gnu/pipewire-0.3/jack
Restart=on-failure
RestartSec=5
# Web UI must not compete with RT audio
Nice=10

[Install]
WantedBy=default.target
