# configs/systemd/user/pi4-audio-webui.service
[Unit]
Description=Pi4 Audio Workstation monitoring web UI (D-020)
After=pipewire.service
Wants=pipewire.service

[Service]
Type=simple
NotifyAccess=main
WorkingDirectory=/home/ela/web-ui
# F-063: Single worker is intentional — multiple workers would break WebSocket
# state sharing (collectors, mode_manager, thermal_monitor are per-process).
# Thread pool expanded to 32 in main.py lifespan to prevent PCM relay blocking.
ExecStart=/home/ela/audio-workstation-venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8080 --workers 1 --ssl-keyfile /home/ela/web-ui/key.pem --ssl-certfile /home/ela/web-ui/cert.pem
Environment=PI_AUDIO_MOCK=0
Environment=PI4AUDIO_SIGGEN=1
Environment=PI4AUDIO_RC_DIR=/home/ela/room-correction
Environment=PI4AUDIO_MEAS_DIR=/home/ela/measurement
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=JACK_NO_START_SERVER=1
Environment=LD_LIBRARY_PATH=/usr/lib/aarch64-linux-gnu/pipewire-0.3/jack
# WatchdogSec disabled — uvicorn startup exceeds 30s on Pi 4
# WatchdogSec=30
Restart=on-failure
RestartSec=2
# Nice=10 removed (F-064): single-worker uvicorn needs full CPU access
# to avoid event loop starvation under WebSocket load
Nice=0

[Install]
WantedBy=default.target
