# configs/systemd/user/pi4audio-graph-manager.service
[Unit]
Description=PipeWire Graph Manager — sole session manager for link topology
# WirePlumber must start first to activate adapter node ports (hardware.audio
# + policy.node profiles). GM depends on ports being present before it can
# create links. WP's auto-linking policies are disabled via
# 90-no-auto-link.conf — it only activates ports, GM handles all linking.
After=pipewire.service wireplumber.service
Requires=pipewire.service
Wants=wireplumber.service

[Service]
Type=simple
ExecStart=%h/bin/pi4audio-graph-manager --mode dj --log-level info
Restart=on-failure
RestartSec=3

# RT scheduling: FIFO/80 for watchdog safety mute deadline guarantee (T-044-4).
# Below PipeWire (FIFO/88) but above all non-audio processes. Ensures the
# watchdog's Mult=0.0 mute path runs within the <21ms response budget even
# under CPU contention.
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=80

# Prevent OOM killer from targeting the session manager during audio
OOMScoreAdjust=-500

# Resource limits for PipeWire clients
LimitMEMLOCK=infinity
LimitRTPRIO=88

# Security hardening — PipeWire client profile (SEC-PW-CLIENT)
# See F-035 for rationale.
NoNewPrivileges=true
PrivateTmp=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
CapabilityBoundingSet=
RestrictSUIDSGID=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged
ProtectSystem=full
ProtectHome=read-only

[Install]
WantedBy=default.target
